# contrib/hama-mce-remoted@.service
[Unit]
Description = HAMA MCE remote control daemon

[Service]
Type=simple
ExecStart = /usr/sbin/hama-mce-remoted --events %I
